# pump_controller.service
# Automatically generated by systemd-sysv-generator

[Unit]
Documentation=man:systemd-sysv-generator(8)
Description=LSB: Heat pump controller
StartLimitInterval=0
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
WorkingDirectory=/opt/pumpcontroller
ExecStart=/opt/pumpcontroller/pump_controller

[Install]
WantedBy=multi-user.target
